# lib/systemd/system/oi-lab-xorg-daemon@.socket
[Unit]
Description=Socket for Xorg server :%i

[Socket]
ListenStream=/tmp/.X11-unix/X%i

[Install]
WantedBy=sockets.target