# lib/systemd/system/oi-lab-autoupgrade.service
[Unit]
Description=Executar atualização do sistema

[Service]
Type=oneshot
ExecStart=/usr/sbin/oi-lab-autoupgrade

