# lib/systemd/system/oi-lab-xorg-daemon@.socket
[Unit]
Description=Socket for Xorg server for SM501 card %i

[Socket]
ListenStream=/tmp/.X11-unix/X9%i

[Install]
WantedBy=sockets.target
