# lib/systemd/system/oi-lab-autoupgrade.timer
[Unit]
Description=Atualizar o sistema automaticamente

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true
 
[Install]
WantedBy=timers.target

