# lib/systemd/system/systemd-poweroff.timer
[Unit]
Description=Desligar o sistema

[Timer]
OnCalendar=*-*-* 17:00:00

[Install]
WantedBy=timers.target

