# lib/udev/rules.d/72-seat-oi-lab-attach-helper.rules
ACTION=="remove", GOTO="seat_end"

SUBSYSTEM=="usb", TAG=="master-of-seat", ENV{ID_SEAT}=="", RUN+="/usr/bin/oi-lab-seat-attach-helper %p"

LABEL="seat_end"
