# debian/oi-lab-proinfo-multi-seat-utils.prerm
#!/bin/sh
set -e

remove() {
    systemctl stop oi-lab-xorg-daemon.service || true
    systemctl stop oi-lab-xorg-daemon.socket || true
    systemctl disable oi-lab-xorg-daemon.socket || true
}

case "$1" in
    remove|deconfigure)
        remove
	;;

    upgrade|failed-upgrade)
        ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
esac
