# lib/systemd/system/display-manager.service.d/oi-lab-multi-seat-config.conf
[Unit]
Wants=oi-lab-multi-seat-config-video.service
After=oi-lab-multi-seat-config-video.service
